# src/config/svc-monitor/contrail-svc-monitor.ini
[program:contrail-svc-monitor]
command=/usr/bin/contrail-svc-monitor --conf_file /etc/contrail/svc_monitor.conf
priority=460
autostart=true
killasgroup=true
stopsignal=KILL
redirect_stderr=true
stdout_logfile=/var/log/contrail/contrail-svc-monitor-stdout.log
stderr_logfile=/dev/null
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)
user=contrail
