# src/discovery/contrail-discovery.ini
[program:contrail-discovery]
command=/usr/bin/contrail-discovery --conf_file /etc/contrail/contrail-discovery.conf --listen_port 911%(process_num)01d --worker_id %(process_num)s
numprocs=1
process_name=%(process_num)s
redirect_stderr=true
stdout_logfile= /var/log/contrail/contrail-discovery-%(process_num)s-stdout.log
stderr_logfile=/dev/null
priority=430
autostart=true
killasgroup=true
stopsignal=KILL
exitcodes=0 ;'expected' exit codes for process (default 0,2)
user=contrail
