# src/discovery/contrail-discovery.ini
[program:contrail-discovery]
command=/usr/bin/contrail-discovery --conf_file /etc/contrail/contrail-discovery.conf
priority=430
autostart=true
killasgroup=true
stopsignal=KILL
redirect_stderr=true
stdout_logfile=/var/log/contrail/contrail-discovery-stdout.log
stderr_logfile=/dev/null
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)
user=contrail
