# src/container/mesos-manager-deprecated/contrail-mesos-manager.ini
[program:contrail-mesos-manager]
command=/usr/bin/contrail-mesos-manager -c /etc/contrail/contrail-mesos.conf
priority=460
autostart=true
killasgroup=true
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/var/log/contrail/contrail-mesos-manager-stdout.log
stderr_logfile=/dev/null
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)
user=contrail
