# src/container/mesos-manager-deprecated/contrail-mesos-manager.initd.supervisord
#!/usr/bin/env bash

# chkconfig: 2345 99 01
# description: Juniper Mesos Network Manager

supervisorctl -s unix:///var/run/supervisord_mesos_manager.sock ${1} `basename ${0}`
