# src/config/vcenter-monitor/java/contrail-vcenter-plugin.ini
[program:contrail-vcenter-plugin]
command=/usr/bin/contrail-vcenter-plugin
priority=420
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/var/log/contrail/contrail-vcenter-plugin-stdout.log
stderr_logfile=/dev/null
user=contrail
