# src/container/kube-manager/contrail-kube-manager.ini
[program:contrail-kube-manager]
command=/usr/bin/contrail-kube-manager -c /etc/contrail/contrail-kubernetes.conf
priority=460
autostart=true
killasgroup=true
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/var/log/contrail/contrail-kube-manager-stdout.log
stderr_logfile=/dev/null
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)
user=contrail
