# src/config/schema-transformer/contrail-schema.ini
[program:contrail-schema]
command=/usr/bin/contrail-schema --conf_file /etc/contrail/schema_transformer.conf
priority=450
autostart=true
autorestart=true
killasgroup=true
stopsignal=KILL
redirect_stderr=true
stdout_logfile=/var/log/contrail/contrail-schema-stdout.log
stderr_logfile=/dev/null
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)
user=contrail
