# src/query_engine/contrail-qe.initd.supervisord
#!/usr/bin/env bash

# chkconfig: 2345 99 01
# description: Juniper Network Virtualization QE
supervisorctl -s unix:///tmp/supervisord_analytics.sock ${1} `basename ${0}`
