# src/config/contrail_issu/contrail-issu-run-sync.ini
[program:contrail-issu-run-sync]
command=/usr/bin/contrail-issu-run-sync --conf_file /etc/contrail/contrail-issu.conf 
numprocs=1
process_name=%(process_num)s
redirect_stderr=true
stdout_logfile= /var/log/contrail/contrail-issu-run-sync-stdout.log
stderr_logfile=/dev/null
priority=440
autostart=true
killasgroup=false
stopsignal=TERM
exitcodes=0
