# src/config/api-server/contrail-api.initd.supervisord
#!/usr/bin/env bash

# chkconfig: 2345 99 01
# description: Juniper Network Virtualization API

processname=`basename ${0}`
for prs in `supervisorctl -s unix:///var/run/supervisord_config.sock status | grep $processname | awk '{ print $1}'`
do
 supervisorctl -s unix:///var/run/supervisord_config.sock ${1} $prs
done

