# src/control-node/contrail-control.ini
[program:contrail-control]
command=/usr/bin/authbind /usr/bin/contrail-control
priority=520
autostart=true
killasgroup=true
stopsignal=KILL
stdout_capture_maxbytes=1MB
redirect_stderr=true
stdout_logfile=/var/log/contrail/contrail-control-stdout.log
stderr_logfile=/dev/null
startsecs=5
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)
user=contrail
