# src/dns/contrail-named.initd.supervisord
#!/usr/bin/env bash

# chkconfig: 2345 99 01
# description: Juniper Network Virtualization Named

supervisorctl -s unix:///tmp/supervisord_control.sock ${1} `basename ${0}`
