# src/dns/contrail-named.ini
[program:contrail-named]
command=/usr/bin/authbind /usr/bin/named -f -c /etc/contrail/dns/named.conf
user=contrail
priority=520
autostart=true
killasgroup=true
stopsignal=KILL
stdout_capture_maxbytes=1MB
startsecs=5
redirect_stderr=true
stdout_logfile=/var/log/contrail/contrail-named-stdout.log
stderr_logfile=/dev/null
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)
user=contrail
