# src/analytics/contrail-snmp-collector/contrail-snmp-collector.ini
[program:contrail-snmp-collector]
command=/usr/bin/contrail-snmp-collector
priority=340
autostart=true
killasgroup=true
stopsignal=KILL
stdout_capture_maxbytes=1MB 
redirect_stderr=true
stdout_logfile=/var/log/contrail/contrail-snmp-collector-stdout.log
stderr_logfile=/var/log/contrail/contrail-snmp-collector-stderr.log
startsecs=5
exitcodes=0                   ; 'expected' exit codes for process (default 0,2)
user=contrail
